# docs/systemd/baby-tracker.service.d/10-tls.conf
[Service]
Environment=BABY_TRACKER_TLS_CERT_PATH=/home/josh/baby-tracker/rpi.tail458584.ts.net.crt
Environment=BABY_TRACKER_TLS_KEY_PATH=/home/josh/baby-tracker/rpi.tail458584.ts.net.key
